# layers/meta-myir-remi/recipes-demos/measy-hmi/qt-demo/pnsn.service
[Unit]
Description=VIN initializing setting
After=systemd-initctl.service

[Service]
ExecStart=/home/eeprom.sh

[Install]
WantedBy=multi-user.target
